# shared/init/vlfdaemon.conf.in
#
# Upstart script for vlfdaemon
#

description "The VLF observatory network daemon"

exec ${CMAKE_INSTALL_PREFIX}/bin/ovldaemon --config-dir=${CMAKE_INSTALL_PREFIX}/share/vlfdaemon

start on runlevel [23]
stop on runlevel [016]

